# containers/config/desktop/autostart/ricing.desktop
[Desktop Entry]
Type=Application
Name=RiceApply
Exec=/usr/local/bin/apply-ricing.sh
X-GNOME-Autostart-enabled=true
NoDisplay=true
